# deploy/php_worker/supervisord.conf
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work --sleep=3 --timeout=60 --tries=3
autostart=true
autorestart=true
numprocs=8
redirect_stderr=true
user=root
stdout_logfile=/var/www/storage/logs/worker.log

[supervisord]
